# root/etc/crontab
SHELL=/bin/sh

* * * * * php /var/www/pelican/artisan schedule:run >> /dev/null 2>&1